# podman/file/Quadlet/file-nginx.container
[Container]
Mount=type=bind,source=/mnt/MainHDD,destination=/HDD
