# podman/certbot/Quadlet/certbot.container
[Container]
Environment=DOMAIN="-d *.3birds.uk"
